# meta-hhctrl/recipes-apps/network-watcher/files/network-watcher.service
[Unit]
Description=Network Watcher

[Service]
ExecStart=/usr/bin/network-watcher

[Install]
WantedBy=multi-user.target