# recipes-apps/hhctrl-mgmt/files/hhctrl-mgmt.service
[Unit]
Description=Henhouse Controller Management Service

[Service]
ExecStart=/usr/bin/hhctrl-mgmt

[Install]
WantedBy=multi-user.target